# cron.daily
#!/usr/bin/bash

YESTERDAY=$(/usr/bin/date --date yesterday '+%Y-%m-%d')