# cron.hourly
#!/usr/bin/bash

# TODAY=$(/usr/bin/date '+%Y-%m-%d')

# /usr/bin/node /home/dgc/databroker/update_thumbnail.js --date $TODAY --device rPI_46_1047_1
